# deploy/ubuntu/synchrohrm.service
[Unit]
Description=SynchroHRM mailbox ingestion service
After=network-online.target
Wants=network-online.target

[Service]
Type=simple
User=__APP_USER__
Group=__APP_GROUP__
WorkingDirectory=__APP_DIR__
Environment=NODE_ENV=production
EnvironmentFile=__APP_ENV_FILE__
ExecStart=__NODE_BIN__ dist/src/index.js
Restart=on-failure
RestartSec=5
NoNewPrivileges=true
PrivateTmp=true
ProtectSystem=full
ProtectHome=true
ReadWritePaths=__APP_DATA_DIR__

[Install]
WantedBy=multi-user.target
